# cns/cns.service
[Unit]
Description=CNS Service
[Service]
ExecStart=/usr/local/cns/cns
[Install]
WantedBy=multi-user.target