# xray/xray.service
[Unit]
Description=Xray Service
[Service]
ExecStart=/usr/local/xray/xray -config /usr/local/xray/xray.json
[Install]
WantedBy=multi-user.target